# PXE/pxe-server-centos7/conf/tftp.service
[Unit]
Description=Tftp Server

[Service]
ExecStart=/usr/sbin/in.tftpd -vvv -s /var/lib/tftpboot
StandardInput=socket
#
# systemctl daemon-reload
#
